# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/etc/init.d/ated_iwpriv
#!/bin/sh /etc/rc.common
#
# Copyright (c) 2019 TP-LINK Technologies Co., Ltd.
# All Rights Reserved.

START=79

start() {
	# TODO: add the calibration process.
	[ "$(/sbin/is_cal)" = "false" ] && {
		echo "wireless is not calibrated, starting ated..." > /dev/console
		/sbin/ated_iwpriv &
		/sbin/ated_iwpriv &
	}

	return 0
}

