# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/etc/init.d/cloud_brd
#!/bin/sh /etc/rc.common

START=98
USE_PROCD=1

start_service() 
{	
	procd_open_instance    	
	procd_set_param command /usr/bin/cloud-brd -c /etc/cloud_config.cfg    	
	procd_set_param respawn    	
	procd_close_instance
}

restart()
{
    rm /tmp/cloud/cloud_status_check_done
    stop
    start
    /etc/init.d/cloud_client restart
}
