# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/etc/init.d/dnsproxy
#!/bin/sh /etc/rc.common
# Copyright (c) 2017 Shenzhen TP-LINK Technologies Co.Ltd.

SERVICE_DAEMONIZE=1

START=84
STOP=84

start() {	
#   disable dnsproxy in AP Mode
	config_load sysmode
	config_get system_mode sysmode mode "Router"
	if [ "$system_mode" != "Router" ]; then
		return 0
	fi

	service_start /usr/bin/dnsproxy -S
	
	return 0
}

stop() {
	service_stop /usr/bin/dnsproxy
	
	return 0
}

restart(){
	stop
	start
}