# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/etc/init.d/time_settings
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2010 OpenWrt.org

START=47

USE_PROCD=1
TS_LIBDIR=/lib/time_settings

TIME_SETTINGS_MOD_ID=279
LOG_INF_SERVICE_START=504
LOG_INF_SERVICE_STOP=505
LOG_INF_SERVICE_RESTART=508

start_service() {
    sys_mode=$(uci get sysmode.sysmode.mode)
    groupinfo="$(cat /tmp/group-info)"
    role=${groupinfo#*role\":\"} && role=${role%%\"*}

    if [ "$sys_mode" == "AP" ] && [ "$role" != "AP" ]; then
		killall -9 update_dst_from_ap
		/usr/sbin/update_dst_from_ap > /dev/null 2>&1 &
    else
        $TS_LIBDIR/set_time	
    fi
	logx -p $$ $TIME_SETTINGS_MOD_ID $LOG_INF_SERVICE_START
}

reload_service() {
	# update dst and reload to /etc/TZ
    groupinfo="$(cat /tmp/group-info)"
    role=${groupinfo#*role\":\"} && role=${role%%\"*}
    role_id=1
    [ "$role" != "AP" ] && role_id=0

    if [ "$role_id" == "0" ]; then
        killall -9 update_dst_from_ap
        /usr/sbin/update_dst_from_ap > /dev/null 2>&1 &
        exit 0
    fi
    dstrule_num=$(uci get systime.@system[0].dstrule_num)
    if [ "$dstrule_num" == "" ]; then
        /usr/sbin/cloud_updateDst > /dev/null 2>&1 &	
    fi
	logx -p $$ $TIME_SETTINGS_MOD_ID $LOG_INF_SERVICE_RESTART
	$TS_LIBDIR/set_time	
}

service_triggers() {
    procd_add_reload_trigger "time_settings"
}
