# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/etc/init.d/uhttpd-fw
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1
PROG=/usr/sbin/uhttpd

start_service() {
    mkdir -p /var/firmwares
	mkdir -p /var/upgrade
	
    procd_open_instance
    procd_set_param command "$PROG" -f -p 127.0.0.1:20003 -D -h /var/firmwares
    procd_set_param respawn
    procd_close_instance
}
