# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/etc/hotplug.d/mode/00-sync
#!/bin/sh

[ "$POINT" = END ] || exit 0

[ "$SRCMODE" = "FAP" -a "$DSTMODE" = "HAP" ] && exit 0

ubus call sync boost >/dev/null
ubus call sync probe >/dev/null

if [ "$SRCMODE" = "RE" ] && [ "$DSTMODE" = "FAP" -o "$DSTMODE" = "HAP" ]; then
    rm -rf /var/run/record_* >/dev/null
    rm -rf /var/run/lookup* >/dev/null
fi
