# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/etc/hotplug.d/mode/50-cloud-reconnect
#!/bin/sh

#echo "========== MODE: POINT = $POINT ===========" > /dev/console

[ "$DSTMODE" = "HAP" ] && exit 0

if [ "$POINT" == "END" -o "$POINT" == "CAPCHANGE" ]; then	
    /etc/init.d/cloud_brd restart
fi
