# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/etc/hotplug.d/mode/99-guest
#!/bin/sh

if [ "$POINT" = "END" -a "$DSTMODE" = "RE" ]; then
	echo "guest in hotplug start to reload @@@@@@@@@@@@@@@@@@@@@@@@@@@@@" > /dev/console
	/etc/init.d/guest_eth reload "bind"
fi


