# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/auto-wifi-network.list
/usr/sbin/awnd
/usr/sbin/handle_plc_bkhl_vlanid
/etc/init.d/awn
/usr/sbin/check_pack
/usr/sbin/online-test-tmp
/usr/sbin/awnweight
