# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/busybox.list
/bin/dmesg
/bin/printenv
/usr/bin/id
/bin/login
/usr/bin/realpath
/bin/rm
/usr/bin/wget
/usr/bin/sha512sum
/bin/sh
/usr/bin/lsof
/sbin/vconfig
/usr/bin/test
/usr/bin/unlzop
/bin/zcat
/bin/fgrep
/bin/ash
/sbin/halt
/sbin/reboot
/bin/true
/sbin/route
/usr/bin/md5sum
/usr/bin/free
/usr/sbin/ntpd
/bin/stty
/usr/bin/time
/usr/bin/passwd
/usr/bin/traceroute
/usr/bin/fold
/bin/vi
/bin/hostname
/usr/bin/tftp
/bin/pwd
/usr/bin/nc
/bin/chown
/etc/init.d/cron
/usr/sbin/addgroup
/bin/mount
/usr/bin/xzcat
/usr/bin/ipcs
/usr/bin/od
/bin/sync
/usr/sbin/adduser
/usr/bin/dirname
/usr/bin/bzip2
/usr/bin/hostid
/sbin/hwclock
/sbin/udhcpc
/usr/bin/whoami
/usr/sbin/inetd
/usr/bin/tac
/bin/ping6
/usr/bin/basename
/usr/bin/logname
/bin/dnsdomainname
/usr/sbin/delgroup
/bin/echo
/usr/bin/less
/usr/sbin/chroot
/usr/bin/unexpand
/sbin/devmem
/usr/sbin/crond
/usr/bin/ftpget
/usr/bin/head
/bin/gunzip
/bin/netstat
/usr/sbin/readprofile
/usr/sbin/ftpd
/bin/busybox
/usr/bin/pwdx
/bin/grep
/usr/bin/clear
/usr/bin/seq
/bin/sleep
/bin/nice
/usr/bin/diff
/usr/bin/groups
/usr/bin/sha1sum
/bin/kill
/bin/catv
/bin/mktemp
/usr/bin/tr
/usr/bin/whois
/bin/mknod
/bin/rmdir
/usr/bin/unzip
/usr/bin/uniq
/usr/bin/tee
/usr/bin/ftpput
/bin/getopt
/usr/bin/uudecode
/usr/bin/bzcat
/usr/bin/hexdump
/etc/init.d/sysntpd
/usr/bin/split
/usr/bin/sha256sum
/bin/date
/usr/bin/wc
/usr/bin/nslookup
/usr/bin/[
/bin/base64
/usr/bin/pgrep
/bin/umount
/bin/ipcalc
/bin/netmsg
/usr/bin/tail
/usr/bin/nmeter
/bin/stat
/usr/bin/pkill
/usr/bin/mkfifo
/etc/crontabs/root
/bin/chgrp
/bin/cp
/usr/bin/mkpasswd
/sbin/poweroff
/usr/bin/[[
/bin/dd
/bin/lock
/usr/bin/expand
/bin/df
/bin/false
/usr/bin/sha3sum
/usr/bin/killall
/usr/sbin/remove-shell
/usr/bin/sort
/bin/mkdir
/sbin/ifconfig
/usr/bin/cmp
/usr/bin/reset
/usr/bin/fuser
/sbin/start-stop-daemon
/usr/bin/pmap
/usr/sbin/dhcprelay
/usr/bin/install
/bin/usleep
/bin/cat
/usr/sbin/arping
/bin/lzop
/etc/init.d/telnet
/usr/bin/lzcat
/usr/bin/smemcap
/usr/bin/telnet
/bin/watch
/usr/sbin/tftpd
/usr/bin/uuencode
/bin/touch
/sbin/mkswap
/bin/tar
/usr/bin/yes
/usr/bin/top
/usr/bin/expr
/sbin/switch_root
/usr/bin/env
/sbin/sysctl
/bin/sed
/usr/bin/ar
/usr/bin/logger
/bin/chmod
/usr/bin/awk
/usr/bin/cksum
/bin/uname
/usr/bin/which
/usr/bin/find
/bin/uncompress
/bin/ln
/usr/bin/pstree
/bin/ls
/usr/bin/cryptpw
/sbin/pivot_root
/usr/bin/dc
/usr/bin/nohup
/usr/bin/cut
/usr/bin/sum
/usr/bin/du
/bin/mv
/usr/bin/crontab
/usr/bin/comm
/usr/bin/printf
/usr/bin/unxz
/usr/bin/ipcrm
/usr/bin/tty
/bin/ping
/usr/bin/timeout
/usr/bin/bunzip2
/usr/sbin/chpasswd
/usr/bin/uptime
/bin/fsync
/usr/sbin/add-shell
/usr/sbin/deluser
/bin/pidof
/bin/egrep
/usr/bin/readlink
/usr/bin/strings
/usr/sbin/udhcpd
/sbin/arp
/usr/bin/patch
/bin/ps
/usr/bin/lzopcat
/bin/cpio
/usr/bin/xargs
/usr/bin/hd
/bin/gzip
/usr/bin/dumpleases
/usr/bin/unlzma
