# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/client_mgmt.list
/usr/bin/client_mgmt_arping
/usr/bin/client_upload
/usr/sbin/manager_role_set
/etc/client-mgmt/MAC_OUI.map
/usr/sbin/open_client_triggers
/usr/bin/client_mgmt
/usr/sbin/get_re_client_list
/etc/init.d/client_mgmt
/etc/client-mgmt/TP_LINK_MAC_OUI.map
