# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/guest-network.list
/etc/hotplug.d/mode/99-guest
/usr/sbin/updateneigh.sh
/usr/sbin/guest-schedule
/etc/init.d/guest-schedule
/etc/init.d/guest_eth
/etc/hotplug.d/mode/99-guest-network
/etc/init.d/guest-bandwidth
/lib/guest-eth/guest-eth.sh
