# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/ipv6-firewall.list
/lib/ipv6-firewall/core_log.sh
/lib/ipv6-firewall/core_firewall.sh
/lib/ipv6-firewall/core.sh
/lib/ipv6-firewall/config.sh
/etc/hotplug.d/mode/99-ipv6-firewall
/etc/init.d/ipv6-firewall
