# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/kmod-shortcut-fe.list
/lib/modules/3.10.108/shortcut-fe.ko
/etc/modules.d/45-shortcut-fe
/lib/modules/3.10.108/shortcut-fe-ipv6.ko
/etc/init.d/shortcut-fe
