# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/kmod-qca-hyfi-bridge.list
/etc/hotplug.d/mode/99-hyfi-bridging
/etc/modules.d/40-qca-hyfi-bridge
/etc/init.d/hyfi-bridging
/lib/modules/3.10.108/hyfi-bridging.ko
