# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/kmod-ipv6-pass-through.list
/lib/netifd/passthrough_reload.sh
/lib/modules/3.10.108/ipv6-pass-through.ko
/lib/netifd/proto/passthrough.sh
/etc/init.d/passthrough.init
