# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/leds.list
/etc/tsched.d/leds_schedule
/usr/lib/libled.so
/usr/sbin/leds_set_nightmode
/etc/init.d/leds_night_mode
/usr/bin/leds
/etc/init.d/leds
