# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/libcloud-api.list
/usr/lib/lua/cloud_lib.so
/usr/bin/cloud-brd
/etc/hotplug.d/iface/98-cloud-brd
/etc/init.d/cloud_brd
/etc/cloud_config.cfg
/usr/lib/lua/cloud_dispatcher.lua
/etc/certificate/2048_newroot.cer
