# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/message-center.list
/usr/sbin/msg_check_inet
/usr/sbin/msg_handle_homecare_alerts
/usr/sbin/update_monthly_msg
/usr/sbin/clean_daily_msg
/usr/sbin/msg_handle_ntp_success_events
/usr/sbin/push_monthly_msg
/etc/tsched.d/script.d/message_center
/etc/tsched.d/script.d/access_client
/etc/init.d/msgcenter
/etc/tsched.d/script.d/iot_device_stats
/usr/bin/msgcenterd
/usr/sbin/push_weekly_msg
/etc/hotplug.d/mode/99-msgcenter
