# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/netifd.list
/lib/netifd/proto/dhcp.sh
/etc/init.d/network
/lib/netifd/netifd-proto.sh
/etc/init.d/networkx
/etc/hotplug.d/iface/00-netstate
/sbin/ifdown
/lib/netifd/dhcp.script
/usr/share/udhcpc/default.script
/lib/network/config.sh
/etc/vlan.d/vlan
/etc/hotplug.d/iface/10-sysctl
/sbin/portset
/sbin/netifd
/sbin/devstatus
/etc/hotplug.d/iface/01-ipv6
/etc/hotplug.d/bridge/01-vlanstate
/etc/hotplug.d/mode/00-network
/etc/hotplug.d/bridge/02-alter-wanif
/etc/hotplug.d/iface/05-stat
/etc/config/portspeed
/sbin/ifstatus
/sbin/ifup
/etc/hotplug.d/net/03-vlanset
/sbin/portspeed
/etc/hotplug.d/iface/80-tcpmss
/sbin/ifmac
