# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/network-roaming.list
/etc/init.d/nrd
/lib/functions/nrd-config.sh
/usr/sbin/start_nrd
/usr/sbin/nrd
/usr/sbin/api_client
/etc/hotplug.d/mode/20-nrd
