# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/parental_control.list
/www/locale/DA_DK.js
/www/locale/FI_FI.js
/usr/lib/lua/luci/model/parental_control.lua
/www/script.js
/usr/sbin/read_url_sig_from_flash
/lib/parental_control/url_lib.sh
/www/pctl_blocking.html
/www/locale/IT_IT.js
/lib/parental_control/core_log.sh
/www/locale/DE_DE.js
/lib/parental_control/core_get_history.sh
/usr/sbin/url-lib-upd
/www/url_lib
/www/locale/ZH_TW.js
/etc/hotplug.d/mode/99-parental-control
/lib/parental_control/config.sh
/www/locale/JP_JP.js
/www/css/style.css
/www/img/security.png
/etc/init.d/parental_control
/lib/parental_control/core.sh
/www/locale/NO_NO.js
/www/locale/ES_ES.js
/sbin/pctl_log_update
/www/locale/FR_FR.js
/lib/parental_control/url_lib
/www/locale/EN_US.js
/www/locale/SV_SE.js
/www/locale/PT_PT.js
/lib/parental_control/core_parental_ctrl.sh
/www/locale/NL_NL.js
/www/favicon.ico
/www/locale/KO_KR.js
/www/bp_language.conf
