# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/qos-v2.list
/etc/hotplug.d/iface/99-qos
/etc/init.d/qos
/sbin/qos_check
/etc/hotplug.d/mode/99-qos
/lib/qos/config.sh
/lib/qos/core.sh
/lib/qos/core_qos.sh
/lib/qos/core_log.sh
