# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/reboot-schedule.list
/etc/tsched.d/reboot_schedule
/usr/sbin/rsched_conf
/etc/init.d/reboot_schedule
/usr/bin/rsched
