# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/sync-server.list
/lib/sync-server/scripts/sync-config
/lib/sync-server/scripts/probe
/lib/sync-server/scripts/detect-re
/usr/lib/lua/sync-script.lua
/lib/sync-server/scripts/force-sync
/etc/hotplug.d/mode/00-sync
/etc/init.d/uhttpd-fw
/usr/lib/lua/update-info.lua
/usr/lib/lua/deco_on_off_iface.lua
/usr/bin/sync-server
/lib/sync-server/scripts/firmwareLte
/usr/lib/lua/comp_whole_leave.lua
/lib/sync-server/scripts/firmware-old
/usr/sbin/mix_network_check
/etc/hotplug.d/ip/00-sync
/etc/init.d/sync-server
/etc/init.d/group-key
/lib/sync-server/scripts/request
/lib/sync-server/scripts/get-info
/usr/sbin/sync-check
/lib/sync-server/scripts/delete-info
/lib/sync-server/scripts/firmware
