# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/time-settings.list
/lib/time_settings/set_time
/etc/init.d/time_settings
/usr/sbin/update_dst_from_ap
/etc/hotplug.d/mode/98-sysntpd
/etc/hotplug.d/iface/98-sysntpd
