# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/time-schedule.list
/usr/lib/libtsched.so
/usr/sbin/tsched
/usr/sbin/tsched_conf
/etc/init.d/tsched
/usr/sbin/save_tmp_config_conf
/etc/tsched.d/save_tmp_config
/usr/lib/lua/tsched_conf.lua
