# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/traffic-stats.list
/etc/init.d/tfstats
/usr/sbin/tfstats
/etc/hotplug.d/mode/30-tfaffic-stats
/etc/config/tfstats
