# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/wanDetect.list
/etc/hotplug.d/iface/07-Internet_detection
/etc/init.d/phyLinkUp
/etc/init.d/wan
/usr/sbin/move_interface
/usr/sbin/speed_diagnose.sh
/etc/hotplug.d/mode/10-dns_correct
/usr/bin/wanDetection
/usr/lib/disconnect_lan.sh
/etc/init.d/phyLinkDown
/usr/bin/isWan
