# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/usr/lib/opkg/info/wifi-config.list
/lib/wifi/wifiupdate_ralink.sh
/lib/wifi/ralinkwifi.sh
/etc/init.d/wifi
/etc/init.d/ssid_steering
/usr/sbin/wifi_led.sh
/sbin/wifi_connect
/lib/wifi/wificommon.sh
/etc/config/ssid-steering
/etc/init.d/blacklist
/sbin/reload_country.sh
/etc/init.d/ated_iwpriv
/sbin/set_cal
/usr/sbin/bridge-control
/sbin/fullband-switch
/usr/sbin/wifimac
/lib/essid/essid.sh
/lib/sbin/wifi_try
/usr/sbin/updateassoclist.sh
/usr/sbin/bridge-hyfi-extend
/usr/sbin/wifi_try
