# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/etc/hotplug.d/iface/00-netstate
[ ifup = "$ACTION" ] && {
	uci_toggle_state network "$INTERFACE" up 1
	uci_toggle_state network "$INTERFACE" connect_time $(sed -ne 's![^0-9].*$!!p' /proc/uptime)
	[ -n "$DEVICE" ] && {
		uci_toggle_state network "$INTERFACE" device "$(uci -q get network.$INTERFACE.ifname)"
		uci_toggle_state network "$INTERFACE" ifname "$DEVICE"
	}
}

if [ ifup = "$ACTION" ] && [ wan = "$INTERFACE" ]; then
	if [ -f /usr/sbin/avirastate_check ] && [ ! -f /tmp/avirastate_checking ]; then
		touch /tmp/avirastate_checking
		/usr/sbin/avirastate_check &
	fi
fi
