# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/etc/hotplug.d/iface/98-sysntpd
#!/bin/sh

if [ "$INTERFACE" == "wan" -a "$ACTION" == "ifup" ]; then
	/etc/init.d/sysntpd restart
	/etc/init.d/time_settings reload
fi

if [ "$INTERFACE" == "wanb"  -a "$ACTION" == "ifup" ]; then
	/etc/init.d/sysntpd restart
fi
