# log/firmware/patool_extraction/M4_4.0_en_1.1.2_Build_20231227_Rel.84005/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_unblob_extracted/M4_4.0_en_1.1.2_Build_20231227_Rel.84005.bin_extract/1984826-13396282.squashfs_v4_le_extract/etc/hotplug.d/net/20-plctestled
#! /bin/sh
#
# 20-plctestled.sh
# Copyright (C) 2019 tpuser <tpuser@liushuaiwei>
#
# Distributed under terms of the MIT license.
#

#plcifname=$(uci get plc.config.PlcIfname)
#if [[ "$(uci get leds.settings.mode)" = "test" ]] && [[ "$INTERFACE" = "$plcifname" ]]; then
#    case $ACTION in
#        remove )
#            ubus call leds switch '{"state":"off"}'
#            ;;
#        add )
#            ubus call leds power '{"state":"on"}'
#            ubus call leds switch '{"state":"on"}'
#            ubus call leds conn '{"state":"online"}'
#            ;;
#    esac
#fi
